# daemon/cmpunlocker.service
[Unit]
Description=CMP 170HX Compute Unlock Daemon
After=systemd-modules-load.service
Before=display-manager.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/cmpunlocker/daemon/watchdog.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
